# tools/cisa-site-assessment/analytics/watcher/supervisor/pipeline-watcher.conf
[program:pipeline-watcher]
command=/opt/psa/venv/bin/python /opt/psa/analytics/watcher/pipeline_watcher.py
directory=/opt/psa/analytics/watcher
user=psa
autostart=true
autorestart=true
startretries=3
startsecs=10
stopwaitsecs=30
stdout_logfile=/var/log/pipeline-watcher/stdout.log
stderr_logfile=/var/log/pipeline-watcher/stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=PYTHONPATH="/opt/psa",PATH="/opt/psa/venv/bin:/usr/local/bin:/usr/bin:/bin"

; Resource limits
priority=999
numprocs=1

